# xcs/.setup/app.supervisor.conf
[program:webapp]
command=/bin/bash -c "cd /app/xcs && nice -n 5 npm run start"
redirect_stderr=true
